# bin/impress
#!/bin/sh

# chkconfig: 2345 90 10
# description: Impress Application Server
# processname: impress

### BEGIN INIT INFO
# Provides:          impress
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Impress Application Server
# Description:       Impress Application Server provudes platform for Node.js applications
### END INIT INFO

usage() {
  echo "service impress {start|stop|restart|status|update}"
  exit 0
}

start() {
  cd /impress
  nohup node --stack-trace-limit=1000 server > /dev/null 2>&1 &
}

stop() {
  killall "impress srv"
}

restart() {
  iprocess=$(ps aux | grep 'impress' | grep -v grep | grep -v sh | wc -l)
  if [ $iprocess -gt 0 ]
  then
    stop
    start
  else
    echo "Impress Application Server is not running"
  fi
}

status() {
  iprocess=$(ps aux | grep 'impress' | grep -v grep | grep -v sh | wc -l)
  if [ $iprocess -gt 0 ]
  then
    ps aux | grep 'impress\|%CPU' | grep -v grep | grep -v sh
  else
    echo "Impress Application Server is not running"
  fi
}

update() {
  cd /impress
  npm update npm -g
  npm update
  restart
}

case $1 in

  start)
    echo "Starting Impress Application Server"
    start
    ;;
  stop)
    echo "Stopping Impress Application Server"
    stop
    ;;
  restart|reload)
    echo "Restarting Impress Application Server"
    restart
    ;;
  update)
    echo "Updating and restarting (if started)"
    update
    ;;
  status)
    status
    ;;
  *) usage
    usage
    ;;

esac
