# bin/impress
#!/bin/sh

# chkconfig: 2345 90 10
# description: Impress Application Server
# processname: impress

### BEGIN INIT INFO
# Provides:          impress
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Impress Application Server
# Description:       Impress Application Server provudes platform for Node.js applications
### END INIT INFO

usage() {
	echo "service impress {start|stop|restart|status|update}"
	exit 0
}

start() {
	cd /impress
	nohup node server > /dev/null 2>&1 &
}

stop() {
	killall "impress srv"
}

status() {
	ps aux | grep impress | grep -v grep | grep -v sh
}

update() {
	cd /impress
	npm update
	stop
	start
}

case $1 in

	start)
		echo -e "\e[1;32mStarting Impress Application Server\e[0m"
		start
		;;
	stop)
		echo -e "\e[1;32mStopping Impress Application Server\e[0m"
		stop
		;;
	restart|reload)
		echo -e "\e[1;32mRestarting Impress Application Server\e[0m"
		stop
		start
		;;
	update)
		echo -e "\e[1;32mUpdating and restarting\e[0m"
		update
		;;
	status)
		status
		;;
	*) usage
		usage
		;;

esac
