# bin/impress
#!/bin/sh

# chkconfig: 2345 90 10
# description: Impress Application Server
# processname: impress

yellow="\e[1;32m"
nocolor="\e[0m"

usage() {
	echo "service impress {start|stop|restart|status}"
	exit 0
}

start() {
	cd /impress
	nohup node server > /dev/null 2>&1 &
}

stop() {
	killall "impress srv"
}

status() {
	ps aux | grep impress | grep -v grep | grep -v sh
}

case $1 in

	start)
		echo -e "${yellow}Starting Impress Application Server${nocolor}"
		start
		;;
	stop)
		echo -e "${yellow}Stopping Impress Application Server${nocolor}"
		stop
		;;
	restart|reload)
		echo -e "${yellow}Restarting Impress Application Server${nocolor}"
		stop
		start
		;;
	status)
		status
		;;
	*) usage
		usage
		;;
esac
